# etc/systemd/system/restart.service
[Unit]
Description=Reboot when there's no internet connection
After=network.target

[Service]
ExecStart=/home/reboot.sh